# mantisstinger/supervisor/httpserv.sv.conf
[program:httpserv]
command=/usr/local/sbin/httpserv -listen=":9000" -root="/www"
redirect_stderr=true