# mantisstinger/supervisor/mantis.sv.conf
[program:mantis]
command=/usr/local/sbin/mantis -config="/conf/mantis.json" -listen=":9001" -privapi=":9002"
redirect_stderr=true