# mantisstinger/supervisor/stinger.sv.conf
[program:stinger]
command=/usr/local/sbin/stinger -config="/conf/stinger.json" -listen=":9003"
redirect_stderr=true